# src/.vuepress/public/scripts/process_blog_access.timer
[Unit]
Description=Run Blog Access Log Processor Daily

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target