# src/.vuepress/public/scripts/process-blog-access.service
[Unit]
Description=Blog Access Log Processor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/var/www/vuepress-blog/scripts/process-blog-access.sh
Restart=always
RestartSec=86400

[Install]
WantedBy=multi-user.target